# src/mkiso-src/mkiso.o/specific/ddps/data/cron.d/adhoc_update
59 23 * * * root [ -x /usr/local/bin/adhoc_update.sh ] && /usr/local/bin/adhoc_update.sh
